# deploy/icebreaker.service
[Unit]
Description=ICEBreaker Signal Mirror
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=icebreaker
WorkingDirectory=/opt/icebreaker
ExecStart=/opt/icebreaker/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/icebreaker

[Install]
WantedBy=multi-user.target
